# launchd/com.awiki.ws-listener.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.awiki.ws-listener</string>

    <key>ProgramArguments</key>
    <array>
        <!-- 替换为实际的 Python 路径（建议使用 venv 中的 python） -->
        <string>/Users/cs/work/agents/awiki-space/awiki-agent-id-skill/.venv/bin/python</string>
        <string>/Users/cs/work/agents/awiki-space/awiki-agent-id-skill/scripts/ws_listener.py</string>
        <string>--credential</string>
        <string>default</string>
        <string>--config</string>
        <string>/Users/cs/work/agents/awiki-space/awiki-agent-id-skill/launchd/listener.example.json</string>
    </array>

    <!-- 开机自动启动 -->
    <key>RunAtLoad</key>
    <true/>

    <!-- 进程退出后自动重启 -->
    <key>KeepAlive</key>
    <true/>

    <!-- 重启最小间隔（秒），防止快速崩溃循环 -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- 工作目录 -->
    <key>WorkingDirectory</key>
    <string>/Users/cs/work/agents/awiki-space/awiki-agent-id-skill</string>

    <!-- 日志输出 -->
    <key>StandardOutPath</key>
    <string>/tmp/awiki-ws-listener.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/awiki-ws-listener.stderr.log</string>

    <!-- 环境变量（按需配置） -->
    <!--
    <key>EnvironmentVariables</key>
    <dict>
        <key>E2E_USER_SERVICE_URL</key>
        <string>https://awiki.ai</string>
        <key>E2E_MOLT_MESSAGE_URL</key>
        <string>https://awiki.ai</string>
        <key>E2E_DID_DOMAIN</key>
        <string>awiki.ai</string>
    </dict>
    -->
</dict>
</plist>
